# upstart_security_camera.conf
start on runlevel 2
#start on startup
respawn limit unlimited
respawn

chdir /home/pi/RasPiMotionDetect

script
  ./RasPiCamera.py 2>&1 >> /var/log/security_camera.log
end script
