# upstart_security_camera.conf
script
  cd /home/pi/RasPiMotionDetect && ./RasPiCamera.py >> /var/log/security_camera.log
end script
start on net-device-up
respawn limit 50
